# server/dcapi.conf
description "My Python - API ZeroMQ"
author "michael.dunham@jusbetter.com"

start on runlevel [2345]
stop on runlevel [!2345]

env AN_ENVIRONMENTAL_VARIABLE=i-want-to-set

respawn

exec /srv/dcapi/server/server.py